# schedule/words-yearly.example.crontab
# NOTE: this won't work with the example wordlist, get a bigger one
5 0 * * * cd /path/to/wordle/wordlist && ./set-word-from.sh words.example.txt "$(date '+%j' | bc)"
